# .config/production/supervisord.conf
[supervisord]
logfile=/var/log/supervisor.log
user=root

[program:nginx]
command=nginx -g "daemon off;"

[program:gunicorn]
command=gunicorn -c /srv/amanda_clone/.config/dev_dooh/gunicorn.py config.wsgi.dev_dooh